# week-summarizer.service
[Unit]
Description=Telegram Weekly Summarizer Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/week-summarizer
ExecStart=/opt/week-summarizer/venv/bin/python scheduler.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
